# roles/ruby/templates/etc/profile.d/ruby.sh.j2
echo $PATH | grep -q "$HOME/.local/share/gem/ruby/{{ ruby__version }}/bin:" || export PATH="$HOME/.local/share/gem/ruby/{{ ruby__version }}/bin:$PATH"
